# article-3/roles/myapp/files/myapp.conf
[program:myapp]
command=java -jar /home/vagrant/myapp-0.0.1-SNAPSHOT.jar
stdout_logfile=/var/log/myapp/myapp.log
stderr_logfile=/var/log/myapp/myapp_error.log
